# Library/LaunchAgents/com.brew.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brew.update</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/brew</string>
      <string>update</string>
    </array>

    <!-- every 7 days (604800 seconds) -->
    <key>StartInterval</key>
    <integer>604800</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Standard output/error logs -->
    <key>StandardOutPath</key>
    <string>/tmp/brew.update.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/brew.update.err.log</string>
</dict>
</plist>
